# ops/com.bytedance.codex.archive-sync-caffeinate-root.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.bytedance.codex.archive-sync-caffeinate.root</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/gongxin/Documents/Projects/strategy-archive-site/scripts/archive-sync-caffeinate.sh</string>
    <string>5400</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>49</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/tmp/com.bytedance.codex.archive-sync-caffeinate.root.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/com.bytedance.codex.archive-sync-caffeinate.root.stderr.log</string>
</dict>
</plist>
